# publications-deb/src/deb/control/prerm
#!/bin/sh
systemctl stop    publications || true
systemctl disable publications || true
